# terraform-gcp/github-runner/scripts/runner-monitor.service
[Unit]
Description=GitHub Runner Health Monitor and Auto-Recovery
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/github-runner/scripts/monitor-runner-health.sh
Restart=always
RestartSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=github-runner-monitor

[Install]
WantedBy=multi-user.target